# server/node.ini
[program:amhub]
command=node /root/AMHub/server/server.js
directory=/root/AMHub
autostart=true
autorestart=true
stderr_logfile=/var/log/supervisor/node.err.log
stdout_logfile=/var/log/supervisor/node.out.log
